# trade_bot/services/postback_service.service
[Unit]
Description=PocketAI Postback Listener
After=network.target

[Service]
User=root
WorkingDirectory=/root/bot
ExecStart=/root/bot/venv/bin/python -m uvicorn postback_app:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
